# deployment/printer-timelapse-sync.timer
[Unit]
Description=Sync timelapse videos from printer every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
